# elastic-network-interfaces.conf
# This task finds and configures elastic network interfaces
# left in an unconfigured state.

start on net-device-up IFACE=eth0
task
exec /sbin/ec2ifscan
